# bot/sovereignlaunch-monitor.service
[Unit]
Description=SovereignLaunch Telegram Monitor
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sovereignlaunch/bot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/ubuntu/sovereignlaunch/bot/.env
ExecStart=/usr/bin/python3 /home/ubuntu/sovereignlaunch/bot/telegram_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sovereignlaunch-monitor

[Install]
WantedBy=multi-user.target
